# db_backups/db_backup.timer
[Unit]
Description=Run database backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target

